# templates/houndd.upstart.erb
# Hound daemon (Upstart unit)
# Version: <%= scope.lookupvar('hound::version') %>
description "Hound daemon"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

env HOUND_BIN=<%= scope.lookupvar('hound::bin_dir') %>/houndd
env HOUND_CONFIG=<%= scope.lookupvar('hound::conf_dir') %>/config.json
env HOUND_ROOT=<%= scope.lookupvar('hound::package_dir') %>

script
exec sudo -u <%= scope.lookupvar('hound::user') %> -g <%= scope.lookupvar('hound::group') %> $HOUND_BIN -conf=$HOUND_CONFIG -root=$HOUND_ROOT -addr=<%= scope.lookupvar('hound::host') %>:<%= scope.lookupvar('hound::port') %>
end script

respawn
respawn limit 10 10
kill timeout 10
